# example/crontab
# Usage
* * * * * root node $WORKDIR/job.js >> /var/log/cron.log 2>&1
